# blog.conf
description "my blog"
author "Costa"
  
start on (local-filesystems and net-device-up)
stop on runlevel [06]
  
respawn
exec /root/website/blog/Costa_Tornado_Blog/server.py